# root/home/.config/systemd/user/radiko_pgm.service
[Unit]
Description = Get Radiko's Program



[Service]
Type = oneshot
ExecStart = /usr/bin/bash %h/.local/bin/sh/radiko_pgm.sh