# root/home/.config/systemd/user/sample.radiko_pgm.service
[Unit]
Description = Get Radiko's Program




[Service]
Type = oneshot
ExecStart = /usr/bin/bash      %h/dotfiles/root/home/bin/bash/radiko_pgm.sh